# supervisord.conf
[supervisord]
nodaemon=true

[program:frontend]
command=npm start

[program:backend]
command=node /app/server/index.js

[program:database]
command=/app/pb/pocketbase serve --http=0.0.0.0:8090
